# ansible/roles/orchestrator/templates/orchestrator.service.j2
[Unit]
Description=Preview Deployer Orchestrator Service
After=network.target docker.service nginx.service
Requires=docker.service nginx.service

[Service]
Type=simple
User={{ deployment_user }}
Group={{ deployment_user }}
WorkingDirectory={{ orchestrator_dir }}
EnvironmentFile={{ orchestrator_dir }}/.env
ExecStart=/bin/bash -c '. ~/.nvm/nvm.sh && exec node {{ orchestrator_dir }}/dist/index.js'
Restart=always
RestartSec=10
StandardOutput=append:{{ orchestrator_dir }}/logs/orchestrator.log
StandardError=append:{{ orchestrator_dir }}/logs/orchestrator-error.log

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
# using true here will cause the service to fail to start, so we use read-only
ProtectHome=read-only
ReadWritePaths={{ orchestrator_dir }} /opt/preview-deployments /etc/nginx/preview-configs

[Install]
WantedBy=multi-user.target
